# electron-app/robot-arm-camera.service
# Systemd unit: USB camera MJPEG stream for the Camera tab.
# install-camera-service.sh copies this to /etc/systemd/system/

[Unit]
Description=Robot Arm USB camera MJPEG stream
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=SERVICE_USER
Group=video
WorkingDirectory=INSTALL_DIR
Environment=ROBOT_ARM_CAMERA_PORT=__CAMERA_PORT__
Environment=ROBOT_ARM_CAMERA_DEVICE=__CAMERA_DEVICE__
ExecStart=INSTALL_DIR/start-camera-stream.sh
StandardOutput=journal
StandardError=journal
Restart=on-failure
RestartSec=5
# Lower scheduling priority so vision processing yields to servo control and kiosk
Nice=10

[Install]
WantedBy=multi-user.target
